# deploy/feaweld-twin.service
[Unit]
Description=feaweld Digital Twin Daemon
Documentation=https://github.com/feaweld/feaweld
After=network.target mosquitto.service

[Service]
Type=notify
ExecStart=/usr/bin/python3 -m feaweld.digital_twin.daemon
Restart=on-failure
RestartSec=10
WatchdogSec=30

# Resource limits
MemoryMax=4G
CPUQuota=200%

# Logging via journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=feaweld-twin

# Security hardening
DynamicUser=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/feaweld

# Environment
Environment=FEAWELD_TMPDIR=/var/lib/feaweld/tmp

[Install]
WantedBy=multi-user.target
